# deploy/openclaw-vllm-small-bench.service
[Unit]
Description=OpenClaw bench Qwen3.5 4B vLLM API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
Environment=CUDA_DEVICE_ORDER=PCI_BUS_ID
Environment=CUDA_VISIBLE_DEVICES=0
Environment=VLLM_NO_USAGE_STATS=1
ExecStart=/home/ubuntu/.venvs/vllm/bin/vllm serve Qwen/Qwen3.5-4B --served-model-name qwen3.5-4b --host 10.68.198.1 --port 8003 --gpu-memory-utilization 0.90 --max-model-len 32768 --dtype auto --api-key vllm-local --reasoning-parser qwen3 --enable-auto-tool-choice --tool-call-parser qwen3_coder --language-model-only --enforce-eager
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
